# roles/plusclouds.os.python.3_8.install/templates/plusclouds-debian.service.j2
[Unit]
Description=Plusclouds Python Service

[Service]
ExecStart=/usr/local/bin/python3.8 -m plusclouds-service
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.Target